# scripts/scheduler/launchd/com.fs.claude-scheduler.watson-puzzle-content.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fs.claude-scheduler.watson-puzzle-content</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fs/Code/game-watson/scripts/scheduler/run-task.sh</string>
        <string>watson-puzzle-content</string>
        <string>/Users/fs/.claude/scheduled-tasks/watson-puzzle-content/SKILL.md</string>
        <string>/Users/fs/Code/game-watson-wt/main-shared</string>
        <string>main-locked</string>
    </array>
    <!-- cron "43 11 * * *" · worktree: main-shared · branch: main · coord: main-locked -->
    <key>StartCalendarInterval</key>
    <dict><key>Minute</key><integer>43</integer><key>Hour</key><integer>11</integer></dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/fs/.claude/external-scheduler/logs/watson-puzzle-content/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fs/.claude/external-scheduler/logs/watson-puzzle-content/launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/fs/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/fs</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>PYTHONIOENCODING</key>
        <string>utf-8</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/fs/Code/game-watson-wt/main-shared</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
